# Supervisor/supervisor/tasks-available/sslocal.ini
[program:sslocal]
user=shadowsocks
command=sslocal -c /etc/shadowsocks.json
process_name=client
autostart=true
autorestart=true
numprocs=1
numprocs_start=3
stdout_logfile=/var/log/shadowsocks/sslocal.log        ; stdout log path, NONE for none; default AUTO
stderr_logfile=/var/log/shadowsocks/sslocal.err        ; stderr log path, NONE for none; default AUTO
